# services/lnk.service
[Unit]
Description=lnk URL shortener

[Service]
User=lnk
Group=lnk
ExecStart=/usr/local/bin/lnk

[Install]
WantedBy=multi-user.target